# systemd/slideshow-server.service
[Unit]
Description=Serve the conference slideshow over localhost:8080
After=network.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=@REPO_DIR@
ExecStart=/usr/bin/python3 @REPO_DIR@/scripts/server.py 8080 --bind 127.0.0.1
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
